# services/battery-monitor.service
[Unit]
Description=Battery Monitor
After=multi-user.target

[Service]
Type=notify
ExecStart=/app/bodycam2/venv/bin/python3 /app/bodycam2/battery/battery_monitor.py
WorkingDirectory=/app/bodycam2/battery

Restart=always
RestartSec=5
WatchdogSec=120
StartLimitIntervalSec=300
StartLimitBurst=10

StandardOutput=null
StandardError=journal

[Install]
WantedBy=multi-user.target
