# services/imu.service
[Unit]
Description=IMU Fall Detection (ICM-42605) Service
After=multi-user.target
StartLimitIntervalSec=0

[Service]
Type=notify
WorkingDirectory=/app/bodycam2/imu/
ExecStart=/app/bodycam2/venv/bin/python3 /app/bodycam2/imu/imu_fall_detect.py
Restart=always
RestartSec=3
WatchdogSec=15
TimeoutStartSec=30
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
